# docker/var/spool/cron/crontabs/root
* * * * * /usr/local/bin/php /var/www/pencatduk-webserver/artisan schedule:run
